# webapp/deploy/gunicorn.service
# systemd unit for the Persian webapp (Gunicorn/WSGI).
# Install: sudo cp deploy/gunicorn.service /etc/systemd/system/persian.service
#          sudo systemctl daemon-reload && sudo systemctl enable --now persian
#
# NOTE (live deploy on linux-surface): cloned to /opt/persian/persian/webapp,
# runs as user 'ashkon', and binds 127.0.0.1:8017 (port 8001 is used by the
# household-tracker app on the same server). Adjust paths/port for your host.

[Unit]
Description=Persian language webapp (Gunicorn)
After=network.target postgresql.service

[Service]
User=ashkon
Group=ashkon
WorkingDirectory=/opt/persian/persian/webapp
EnvironmentFile=/opt/persian/persian/webapp/.env
ExecStart=/opt/persian/persian/webapp/.venv/bin/gunicorn config.wsgi:application \
    --bind 127.0.0.1:8017 \
    --workers 3 \
    --timeout 60
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
